# systemd/i3-battery-nagbar.service
[Unit]
Description=battery nagbar service

[Service]
ExecStart=/usr/bin/i3-battery-nagbar

[Install]
WantedBy=basic.target
